# raspberry-pi/systemd/eden-hangwa.service
[Unit]
Description=Eden Hangwa Order Management System
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/eden-hangwa
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

# 환경변수 (PostgreSQL 사용)
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512
Environment=PORT=3000
# DATABASE_URL은 .env 파일에서 로드됩니다
EnvironmentFile=/home/pi/eden-hangwa/.env

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eden-hangwa

[Install]
WantedBy=multi-user.target
